# ivi/99-atkbd.rules
KERNEL=="id", SUBSYSTEM=="dmi", ATTR{sys_vendor}=="VMware, Inc.", RUN+="/sbin/modprobe atkbd"
